# zfs-auto-snapshot.timer.j2
[Unit]
Description=ZFS snapshot timer

[Timer]
OnCalendar={{ zfs_snapshot_schedule }}
Persistent=true

[Install]
WantedBy=timers.target